# deploy/gunicorn.service
[Unit]
Description=MIBSP Gunicorn Application Server
After=network.target mysql.service redis.service

[Service]
User=www-data
Group=www-data
WorkingDirectory=/var/www/mibsp
Environment="PATH=/var/www/mibsp/venv/bin"
EnvironmentFile=/var/www/mibsp/.env
ExecStart=/var/www/mibsp/venv/bin/gunicorn \
    --workers 4 \
    --bind 127.0.0.1:8000 \
    --timeout 120 \
    --access-logfile /var/log/mibsp/access.log \
    --error-logfile /var/log/mibsp/error.log \
    --capture-output \
    --enable-stdio-inheritance \
    wsgi:app
ExecReload=/bin/kill -s HUP $MAINPID
KillMode=mixed
TimeoutStopSec=5
PrivateTmp=true
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
